# deploy/realm.service
[Unit]
Description=Realm Game Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/realm/server
Environment="PATH=/home/ubuntu/realm/server/venv/bin"
ExecStart=/home/ubuntu/realm/server/venv/bin/gunicorn \
    --worker-class eventlet \
    -w 1 \
    --bind 127.0.0.1:5000 \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
